# ansible/templates/orderbook-collector.service.j2
# Systemd service template for Orderbook Collector
# Quants-Lab Data Collection Service

[Unit]
Description=Orderbook Collector - {{ exchange | upper }}
Documentation=https://github.com/yourorg/quants-lab
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User={{ service_user | default('ubuntu') }}
Group={{ service_group | default('ubuntu') }}
WorkingDirectory={{ work_dir | default('/opt/quants-lab') }}

# 环境变量
Environment="PATH={{ conda_path | default('/home/ubuntu/miniconda3/envs/quants-lab/bin') }}:/usr/local/bin:/usr/bin:/bin"
Environment="EXCHANGE={{ exchange }}"
Environment="METRICS_PORT={{ metrics_port | default(8001) }}"
Environment="LOG_LEVEL={{ log_level | default('INFO') }}"
Environment="PYTHONUNBUFFERED=1"

# 启动命令
ExecStart={{ python_path | default('/home/ubuntu/miniconda3/envs/quants-lab/bin/python') }} \
    cli.py run-tasks \
    --config {{ config_file | default('config/orderbook_snapshot_' + exchange + '.yml') }}

# 停止命令
ExecStop=/bin/kill -SIGTERM $MAINPID

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# 资源限制
MemoryLimit={{ memory_limit | default('2G') }}
CPUQuota={{ cpu_quota | default('100%') }}

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orderbook-collector-{{ exchange }}

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

